# ttserve.conf
description "ttserve"
start on startup
task
exec /home/ubuntu/go/src/github.com/rayozzie/teletype-ttserve/upstart.sh
